# infra/launchd/com.news.enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.news.enrich</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>${REMOTE_DIR}/infra/scripts/enrich.py</string>
  </array>
  <key>WorkingDirectory</key><string>${REMOTE_DIR}</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MOONSHOT_API_KEY</key><string>${MOONSHOT_API_KEY}</string>
    <key>MOONSHOT_MODEL</key><string>${MOONSHOT_MODEL}</string>
    <key>MOONSHOT_API_URL</key><string>${MOONSHOT_API_URL}</string>
    <key>PGDATABASE</key><string>${PGDATABASE}</string>
    <key>PGUSER</key><string>${PGUSER}</string>
    <key>PGHOST</key><string>${PGHOST}</string>
    <key>PGPORT</key><string>${PGPORT}</string>
    <key>PGPASSWORD</key><string>${PGPASSWORD}</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>${REMOTE_USER_HOME}/Library/Logs/news-enrich.log</string>
  <key>StandardErrorPath</key><string>${REMOTE_USER_HOME}/Library/Logs/news-enrich.log</string>
</dict>
</plist>
