# infra/launchd/com.news.rsshub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.news.rsshub</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/coco/code/rsshub/dist/index.mjs</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/coco/code/rsshub</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key><string>production</string>
    <key>CACHE_TYPE</key><string>memory</string>
    <key>PORT</key><string>1200</string>
    <key>NODE_OPTIONS</key><string>--max-http-header-size=32768</string>
    <key>TWITTER_AUTH_TOKEN</key><string></string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/coco/Library/Logs/news-rsshub.log</string>
  <key>StandardErrorPath</key><string>/Users/coco/Library/Logs/news-rsshub.log</string>
</dict>
</plist>
