# infra/launchd/com.news.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.news.ingest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>${REMOTE_DIR}/infra/scripts/manual-ingest.py</string>
  </array>
  <key>WorkingDirectory</key><string>${REMOTE_DIR}</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PROXY</key><string>${MBP_PROXY}</string>
    <key>RSSHUB</key><string>http://127.0.0.1:${RSSHUB_PORT}</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>${REMOTE_USER_HOME}/Library/Logs/news-ingest.log</string>
  <key>StandardErrorPath</key><string>${REMOTE_USER_HOME}/Library/Logs/news-ingest.log</string>
</dict>
</plist>
