# infra/launchd/com.news.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.news.web</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>${REMOTE_DIR}/web/node_modules/.bin/next</string>
    <string>start</string>
  </array>
  <key>WorkingDirectory</key><string>${REMOTE_DIR}/web</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key><string>production</string>
    <key>N8N_WEBHOOK_BASE</key><string>${N8N_WEBHOOK_BASE}</string>
    <key>PORT</key><string>${WEB_PORT}</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>${REMOTE_USER_HOME}/Library/Logs/news-web.log</string>
  <key>StandardErrorPath</key><string>${REMOTE_USER_HOME}/Library/Logs/news-web.log</string>
</dict>
</plist>
