# infra/launchd/com.news.enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.news.enrich</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/coco/code/news/infra/scripts/enrich.py</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/coco/code/news</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MOONSHOT_API_KEY</key><string>__FILL_FROM_ENV__</string>
    <key>MOONSHOT_MODEL</key><string>kimi-k2-0905-preview</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/coco/Library/Logs/news-enrich.log</string>
  <key>StandardErrorPath</key><string>/Users/coco/Library/Logs/news-enrich.log</string>
</dict>
</plist>
